# supervisor.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
environment=APP_HOME="%(ENV_APP_HOME)s",MCPJUNGLE_DATA_ROOT="%(ENV_MCPJUNGLE_DATA_ROOT)s",HOME="%(ENV_HOME)s",PATH="%(ENV_PATH)s",LANG="%(ENV_LANG)s",LC_ALL="%(ENV_LC_ALL)s",TMPDIR="%(ENV_TMPDIR)s",XDG_CONFIG_HOME="%(ENV_XDG_CONFIG_HOME)s",XDG_CACHE_HOME="%(ENV_XDG_CACHE_HOME)s",XDG_DATA_HOME="%(ENV_XDG_DATA_HOME)s",PYTHONPATH="%(ENV_PYTHONPATH)s"

[program:nginx]
command=nginx -c /app/code/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcpjungle]
command=/usr/local/bin/mcpjungle start
user=cloudron
environment=PORT="8081"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:admin-api]
command=python3 -m mcpjungle_admin.admin_api
user=cloudron
environment=MCPJUNGLE_ADMIN_TOKEN="%(ENV_MCPJUNGLE_ADMIN_TOKEN)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:boot-reconcile]
command=/bin/bash -c "sleep 5 && /usr/local/bin/mcpjungle init-server --registry http://127.0.0.1:8081 2>&1 || true && /usr/local/bin/mcpjungle-admin reconcile --boot-mode && /usr/local/bin/mcpjungle-admin sync-groups && /usr/local/bin/mcpjungle-admin generate-lazy-config"
user=cloudron
autostart=true
autorestart=false
startsecs=0
exitcodes=0,1
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
